# config/systemd/user/bg-rotate.service
[Unit]
Description=Automatically rotate background images

[Service]
PassEnvironment=PATH
ExecSearchPath=%h/scripts
ExecStart=bg-rotate --all

# Hardening
SystemCallArchitectures=native
MemoryDenyWriteExecute=true
NoNewPrivileges=true

# Elevated permissions to sync ownership (disabled by default),
# see https://docs.syncthing.net/advanced/folder-sync-ownership
#AmbientCapabilities=CAP_CHOWN CAP_FOWNER


